# crontab/diskusage.sh
#! /bin/bash

# printing the current date
date

# printing disk statistics
df -h
