# crontab/everymin.sh
#! /bin/bash

# prints the date every minute
date >> /tmp/everymin.txt

